# etc/systemd/system/isu-go.service.d/runtime.conf
[Service]
Environment=GOMEMLIMIT=2GiB
Environment=GOGC=200
